# SmartHome_v2/ControlPlane/virtual_dataplane/light/supervisord.conf
[unix_http_server]
file=/dev/shm/supervisor.sock
chmod=0755  

[supervisord]
nodaemon=true
pidfile=/tmp/supervisord.pid 

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:8001        ; (ip_address:port specifier, *:port for all iface)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///dev/shm/supervisor.sock

[program:lightmanager]
directory=/code
command=python -m light.light_manager
autostart=true
autorestart=true
stdout_logfile=/code/light.log
stderr_logfile=/code/light.err
